# pbin/android/n4/crontab
# nexus4 crontab
0 4 * * * /sdcard/nosync/profile/pbin/android/n4/n4upkeep.sh